# scripts/se.mazeworks.trygghetsnod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>se.mazeworks.trygghetsnod</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/frej/Code/mazeworks/krisnod/scripts/start.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>KOMMUN</key>
        <string>arvika</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/frej/Code/mazeworks/krisnod/storage/logs/launchagent.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/frej/Code/mazeworks/krisnod/storage/logs/launchagent.err</string>
</dict>
</plist>
